# xiotech/Wookiee/Linux/Nightly/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.30083 installed on Mon Jul 13 19:56:23 2009)
# (Cron version -- $Id: crontab 101320 2009-09-11 12:35:54Z eng_build $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (Cron version -- $Id: crontab 101320 2009-09-11 12:35:54Z eng_build $)
#
# Environment variable setting:    Name = Value.
SHELL = /bin/sh
MAILTO = Marshall_Midden@xiotech.com
PATH=/usr/bin:/bin:/usr/local/bin

# minute hour day_of_month month day_of_week
30 18 * * *	root su -l eng_build -c /home/eng_build/Nite/Nightly/Nite 2>&1 | tee /home/eng_build/Nite/Nightly/lastrun.out |  /usr/bin/mail -s "NIGHTLY `/bin/date` RUN" Marshall_Midden@xiotech.com,Neal_Eiden@xiotech.com
