# xiotech/Wookiee/Linux/xio3d/files-4700/etc/init.d/boot.clock
#! /bin/sh
#
# Copyright (c) 2001-2002 SuSE Linux AG, Nuernberg, Germany.
# All rights reserved.
#
# /etc/init.d/boot.clock
#
### BEGIN INIT INFO
# Provides:          boot.clock
# Required-Start:    boot.localfs
# X-UnitedLinux-Should-Start: boot.crypto boot.swap
# Required-Stop:
# Default-Start:     B S
# Default-Stop:
# Description:       set cmos clock
### END INIT INFO

. /etc/rc.status
. /etc/sysconfig/clock
rc_reset

case "$1" in
  start|restart)
    # set and adjust the CMOS clock
    echo -n Setting up the CMOS clock
    rm -f /etc/adjtime
    test -f /etc/adjtime || echo -e "0.0 0 0.0\n0\nUTC" > /etc/adjtime
    /sbin/hwclock --utc --adjust $HWCLOCK
    rc_status
    /sbin/hwclock --utc --hctosys $HWCLOCK
    rc_status
    rc_status -v -r
    ;;
  stop)
    echo -n "Set Hardware Clock to the current System Time"
    # write back to hardware clock and calculate adjtime
    /sbin/hwclock --utc --systohc $HWCLOCK
    rc_status -v -r
    ;;
  status)
    rc_failed 4
    rc_status -v
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart}"
    exit 1
    ;;
esac

rc_exit

# vi:sw=8 ts=8 noexpandtab
